# services/tempo/run
#!/bin/sh
# ╭――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――╮
# │ TEMPO - S6 SERVICE SCRIPT                                             │
# ╰――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――――╯
# Manages the Tempo process under s6 supervision. Ensures runtime directories
# exist with the correct ownership before launching the binary.

set -euo pipefail

APP_USER="tempo"
CONFIG_FILE="${TEMPO_CONFIG:-/etc/tempo/tempo.yaml}"
DATA_PATH="/var/tempo"

mkdir -p "${DATA_PATH}"
chown -R "${APP_USER}:${APP_USER}" "${DATA_PATH}"

exec s6-setuidgid "${APP_USER}" \
  /usr/bin/tempo \
    -config.file="${CONFIG_FILE}"
